# scripts/installer/flywall-ctl.init
#!/sbin/openrc-run

name="flywall-ctl"
description="Flywall Firewall Control Plane"
command="/usr/local/bin/flywall"
command_args="ctl /etc/flywall/flywall.hcl"
command_background=true
pidfile="/run/flywall-ctl.pid"

depend() {
    need localmount
    before net
    provide firewall
}
